# scripts/juliana-printer-server.service
[Unit]
Description=Juliana Printer Server (CUPS Bridge)
After=network.target juliana-bluetooth.service
Wants=network.target juliana-bluetooth.service

[Service]
Type=simple
User=tato
WorkingDirectory=/home/tato/juliana-orderflow
EnvironmentFile=-/etc/default/juliana-printer-server
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=PRINTER_80MM_NAME=GLPrinter_80mm
Environment=PRINTER_58MM_NAME=PT210A_58mm
ExecStart=/usr/bin/env node /home/tato/juliana-orderflow/printer-server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
